# misc/discord-linux-rpc.service.in
[Unit]
Description=Discord rich presence process monitor

[Service]
Exec=@prefix@/@bindir@/discord-linux-rpc

[Install]
WantedBy=default.target
